# config/systemd/ruos-agent.service
[Unit]
Description=ruos-agent — agentic heartbeat (observe-decide-act loop)
After=ruvultra-brain.service ruvultra-embedder.service
Wants=ruvultra-brain.service ruvultra-embedder.service

[Service]
Type=oneshot
ExecStart=/home/ruvultra/.local/bin/ruos-agent
StandardOutput=journal
StandardError=journal
TimeoutStartSec=120
