# config/systemd/ruos-update.timer
[Unit]
Description=ruos-update — weekly OTA check (Sundays 4 AM)

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
